# aranet4-monitor.service
[Unit]
Description=Aranet4 Metrics Exporter
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/tmp
ExecStart=/usr/bin/python3 /home/daniel/home-hub/aranet4-monitor.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Bluetooth access
SupplementaryGroups=bluetooth

[Install]
WantedBy=default.target
